# servidor/host-agent/supabase-host-agent.service
# Unit template do host-agent. O install.sh substitui os placeholders pelos
# caminhos reais e instala em /etc/systemd/system/.
[Unit]
Description=Supabase multitenant host-agent (comandos fechados de lifecycle)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=exec
User=__HOST_AGENT_USER__
Environment="HOST_AGENT_ROOT=__SERVIDOR_DIR__"
ExecStartPre="__AGENT_DIR__/.venv/bin/python" -m hostagent --root "__SERVIDOR_DIR__" --wait-for-schema
ExecStart="__AGENT_DIR__/.venv/bin/python" -m hostagent --root "__SERVIDOR_DIR__"
WorkingDirectory=__AGENT_DIR__
Restart=always
RestartSec=5
# A espera tem timeout proprio (HOST_AGENT_SCHEMA_WAIT_TIMEOUT).
TimeoutStartSec=0
# Comandos longos (rename/duplicate) precisam terminar ou fazer rollback.
TimeoutStopSec=600
KillMode=mixed
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
